# templates/example.service.j2
[Unit]
Description=Our example of go application
After=network.target

[Service]
Type=simple
ExecStart=/home/vagrant/example1
PIDFile=/var/run/simple

[Install]
WantedBy=default.target
